# JDFruitAndPet/rootfs/etc/services.d/JDFruitAndPet/run
#!/usr/bin/with-contenv bashio

declare name
declare value

bashio::log.info "Initializing..."

cd /;./modifyCronFile.sh

exec crond -f
